# stow/clipster/.config/systemd/user/clipster.service
[Unit]
Description=Clipster clipboard manager daemon

[Service]
ExecStart=clipster --daemon
