# stow/redshift/.config/systemd/user/geoclue-demo.service
[Unit]
Description=Geoclue2 demo agent

[Service]
ExecStart=/usr/lib/geoclue-2.0/demos/agent

[Install]
WantedBy=redshift.service
